# supervisor/udp_server.conf
[program:udp_server]
command=php /var/www/udp_server.php
numprocs=1
directory=/tmp
autostart=true
autorestart=true
startsecs=5
startretries=10
redirect_stderr=false
stdout_logfile=/var/www/logs/udp_server.out.log
stdout_capture_maxbytes=1MB
stderr_logfile=/var/www/logs/udd_server.error.log
stderr_capture_maxbytes=1MB